# playbooks/roles/web_admin/templates/web_admin_d
#!/bin/bash
# /etc/init.d/web_admin_d

### BEGIN INIT INFO
# Provides:          web_admin_d
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: web admin demon
# Description:       used to manage server 
### END INIT INFO

LOGFILE="log/output.txt"

divider () {
  echo "---------------------------------------------" >> $LOGFILE
}

header () {
  divider
  echo $1 >> $LOGFILE
  date >> $LOGFILE
  divider
}

set_env () {
  export GEM_PATH=/root/.gem/ruby/1.9.1:/var/lib/gems/1.9.1
  cd /home/{{deploy_user}}/apps/web_admin
  mkdir -p log
}

show_env () {
  date >> $LOGFILE
  whoami >> $LOGFILE
  echo $PATH >> $LOGFILE
  which ruby >> $LOGFILE
  which thin >> $LOGFILE
  which gem  >> $LOGFILE
  gem list   >> $LOGFILE
  divider
}

start_service () {
  header "Starting web_admin_d" 
  set_env
  show_env
  RACK_ENV=production bin/run > $LOGFILE 2>&1 &
}

stop_service () {
  header "Stopping web_admin_d" 
  pkill -f web_admin_d
}

case "$1" in 
  start) start_service
    ;;
  stop) stop_service
    ;;
  restart)
    stop_service
    sleep 1
    start_service
    ;;
  *)
    echo "Usage: /etc/init.d/web_admin_d start|stop"
    exit 1
    ;;
esac

exit 0

